# ansible/playbooks/templates/face-api.service.j2
[Unit]
Description=Face API Application
Documentation=https://github.com/YOUR_USERNAME/whos-your-papa-ai
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ app_dir }}
ExecStart=/usr/bin/docker-compose -f {{ app_dir }}/docker-compose.yml up -d
ExecStop=/usr/bin/docker-compose -f {{ app_dir }}/docker-compose.yml down
ExecReload=/usr/bin/docker-compose -f {{ app_dir }}/docker-compose.yml restart
TimeoutStartSec=300
TimeoutStopSec=300
User={{ app_user }}
Group={{ app_user }}

# Restart policy
Restart=on-failure
RestartSec=30s

# Environment
Environment=COMPOSE_PROJECT_NAME={{ project_name }}

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths={{ app_dir }}

[Install]
WantedBy=multi-user.target